# kzl-linux/pacman/pkgbuilds/kzl/shadow/shadow.install
post_install() {
    systemctl enable shadow.timer
}
